# init/nginx-chrooter.conf
description     "mount nginx-chrooter"

start on startup

exec if [[ -e /opt/nginx-chrooter-scripts ]] then; run-parts --regex=".*sh" /opt/nginx-chrooter-scripts; fi;

