# init/nginx-chrooter.conf
description     "mount nginx-chrooter"

start on startup

script
    if [[ -e /opt/nginx-chrooter-scripts ]]; then
	run-parts /opt/nginx-chrooter-scripts
    fi
end script